# runit/netbox-dnsmasq-dhcp/run
#!/bin/sh

while true ; do
    /usr/local/bin/netbox-dhcp-hosts > /run/dhcp-hosts && pkill -SIGHUP dnsmasq
    sleep 600
done
